# systemd-user/.config/systemd/user/monitor-disk.service
[Unit]
Description=Disk monitor (notify when partitions >85% full)

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/monitor-disk.sh
