# systemd-user/.config/systemd/user/monitor-battery.service
[Unit]
Description=Battery monitor (notify on low/critical)

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/monitor-battery.sh
